# templates/kodi_bg_changer.timer.jinja2
[Unit]
Description=Change the Kodi background image every hour

[Timer]
Unit=kodi_bg_changer
OnCalendar=hourly
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
